# hello_world_addon/rootfs/etc/services.d/hello/run
#!/usr/bin/with-contenv bash
exec echo "Hello from addon"
