# foxess_inverter/rootfs/etc/services.d/foxess_t_series_integration/run
#!/usr/bin/with-contenv bash
exec python3 /foxess_t_series_integration.py